# ops/systemd/grouped-runner-saturday.timer
[Unit]
Description=Start grouped runner on Saturday at 08:00 America/Sao_Paulo

[Timer]
OnCalendar=Sat 08:00
AccuracySec=1min
Unit=grouped-runner.service

[Install]
WantedBy=timers.target
