# deploy/systemd/overseas-grouped-runner-saturday.service
[Unit]
Description=Overseas grouped runner (saturday)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=otavio
WorkingDirectory=/home/otavio/overseas-business-online
Environment=TZ=America/Sao_Paulo
ExecStart=/home/otavio/overseas-business-online/scripts/systemd_start_grouped_runner.sh 420
ExecStop=/home/otavio/overseas-business-online/scripts/systemd_stop_and_cleanup.sh
KillMode=control-group
TimeoutStopSec=120
Restart=no
StandardOutput=append:/home/otavio/overseas-business-online/.streamlit_logs/grouped_runner_systemd.log
StandardError=append:/home/otavio/overseas-business-online/.streamlit_logs/grouped_runner_systemd.log

[Install]
WantedBy=multi-user.target
